# dotfiles/ai/services/docker-registry.service
[Unit]
Description=Docker Registry
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/registry serve /etc/docker/registry/config.yml
Restart=always
Environment=REGISTRY_STORAGE_FILESYSTEM_ROOTDIRECTORY=/var/lib/registry

[Install]
WantedBy=multi-user.target
